# services/hs_torrent_search.target
[Unit]
Description=HS Torrent Downloader

Requires=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target
